# systemd/trash-man.service
[Unit]
Description=Cleanup my Trash directory
After=default.target

[Service]
ExecStart=trash_man

[Install]
WantedBy=multi-user.target
